# usr/src/cmd/krb5/kwarn/ktkt_warn.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!--
	Copyright 2015 Nexenta Systems, Inc.  All rights reserved.

	Copyright 2005 Sun Microsystems, Inc.  All rights reserved.
	Use is subject to license terms.

	NOTE:  This service manifest is not editable; its contents will
	be overwritten by package or patch operations, including
	operating system upgrade.  Make customizations in a different
	file.

	Service manifest for the Kerberos V5 warning messages daemon
-->

<service_bundle type='manifest' name='SUNWkrbr:kwarn'>

<service
	name='network/security/ktkt_warn'
	type='service'
	version='1'>

	<create_default_instance enabled='false' />

	<restarter>
		<service_fmri value='svc:/network/inetd:default' />
	</restarter>

	<dependency name='rpcbind'
	    grouping='require_all'
	    restart_on='restart'
	    type='service'>
		<service_fmri value='svc:/network/rpc/bind' />
	</dependency>

	<exec_method
		type='method'
		name='inetd_start'
		exec='/usr/lib/krb5/ktkt_warnd'
		timeout_seconds='0'>
		<method_context>
			<method_credential
				user='root'
				group='other'
				privileges='basic,!file_link_any,!proc_info,!proc_session,proc_setid,file_chown,file_dac_read,file_dac_write'
			/>
		</method_context>
	</exec_method>

	<exec_method
		type='method'
		name='inetd_offline'
		exec=':kill_process'
		timeout_seconds='0'>
	</exec_method>

	<exec_method
		type='method'
		name='inetd_disable'
		exec=':kill'
		timeout_seconds='0'>
	</exec_method>

	<property_group name='inetd' type='framework'>
		<propval name='name' type='astring' value='100134' />
		<propval name='endpoint_type' type='astring' value='tli' />
		<propval name='proto' type='astring' value='ticotsord' />
		<propval name='wait' type='boolean' value='true' />
		<propval name='isrpc' type='boolean' value='true' />
		<propval name='rpc_low_version' type='integer' value='1' />
		<propval name='rpc_high_version' type='integer' value='1' />
	</property_group>

	<stability value='Unstable' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
			Kerberos V5 warning messages daemon
			</loctext>
		</common_name>
		<documentation>
			<manpage title='ktkt_warnd' section='1M'
			    manpath='/usr/share/man' />
		</documentation>
	</template>
</service>

</service_bundle>
