# usr/src/cmd/zoneadm/resource-mgmt.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
 Copyright 2006 Sun Microsystems, Inc.  All rights reserved.
 Use is subject to license terms.

 CDDL HEADER START

 The contents of this file are subject to the terms of the
 Common Development and Distribution License (the "License").
 You may not use this file except in compliance with the License.

 You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
 or http://www.opensolaris.org/os/licensing.
 See the License for the specific language governing permissions
 and limitations under the License.

 When distributing Covered Code, include this CDDL HEADER in each
 file and include the License file at usr/src/OPENSOLARIS.LICENSE.
 If applicable, add the following below this CDDL HEADER, with the
 fields enclosed by brackets "[]" replaced with your own identifying
 information: Portions Copyright [yyyy] [name of copyright owner]

 CDDL HEADER END

    NOTE:  This service manifest is not editable; its contents will
    be overwritten by package or patch operations, including
    operating system upgrade.  Make customizations in a different
    file.
-->

<service_bundle type='manifest' name='SUNWzoner:zones'>

<!--
	This service applies global zone resource management settings
	at system startup.
-->
<service
	name='system/resource-mgmt'
	type='service'
	version='1'>

	<create_default_instance enabled='true' />

	<single_instance />

	<dependency
		name='usr'
		type='service'
		grouping='require_all'
		restart_on='none'>
		<service_fmri value='svc:/system/filesystem/minimal' />
	</dependency>

	<dependency
		name='scheduler'
		type='service'
		grouping='optional_all'
		restart_on='none'>
		<service_fmri value='svc:/system/scheduler' />
	</dependency>

	<dependency
		name='pools'
		type='service'
		grouping='optional_all'
		restart_on='none'>
		<service_fmri value='svc:/system/pools' />
	</dependency>

	<dependent
		name='rcap'
		grouping='optional_all'
		restart_on='none'>
		<service_fmri value='svc:/system/rcap' />
	</dependent>

	<exec_method
		type='method'
		name='start'
		exec='/lib/svc/method/svc-resource-mgmt %m'
		timeout_seconds='60'>
	</exec_method>

	<exec_method
		type='method'
		name='stop'
		exec=':true'
		timeout_seconds='3'>
	</exec_method>

	<property_group name='startd' type='framework'>
		<propval name='duration' type='astring' value='transient' />
	</property_group>

	<stability value='Unstable' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
			Global zone resource management settings
			</loctext>
		</common_name>
		<documentation>
			<manpage title='zones' section='7' manpath='/usr/share/man' />
			<manpage
				title='zonecfg'
				section='8'
				manpath='/usr/share/man' />
		</documentation>
	</template>
</service>

</service_bundle>
