# usr/src/cmd/krb5/krb5kdc/krb5kdc.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
	Copyright 2008 Sun Microsystems, Inc.  All rights reserved.
	Use is subject to license terms.

	NOTE:  This service manifest is not editable; its contents will
	be overwritten by package or patch operations, including
	operating system upgrade.  Make customizations in a different
	file.

	Service manifest for the Kerberos key distribution center
-->

<service_bundle type='manifest' name='SUNWkrbr:krb5kdc'>

<service
	name='network/security/krb5kdc'
	type='service'
	version='1'>

	<create_default_instance enabled='false' />

	<single_instance/>

	<dependency
		name='dns'
		grouping='require_all'
		restart_on='error'
		type='service'>
		<service_fmri value='svc:/network/dns/client' />
	</dependency>

	<exec_method
		type='method'
		name='start'
		exec='/usr/lib/krb5/krb5kdc'
		timeout_seconds='60'>
		<method_context working_directory='/'>
			<method_credential
				user='root'
				group='root'
				privileges='basic,!file_link_any,!proc_info,!proc_session,net_privaddr,proc_audit'
			/>
		</method_context>
	</exec_method>

	<exec_method
		type='method'
		name='stop'
		exec=':kill'
		timeout_seconds='60' />

	<stability value='Unstable' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
			Kerberos key distribution center
			</loctext>
		</common_name>
		<documentation>
			<manpage title='krb5kdc' section='8'
			    manpath='/usr/share/man' />
		</documentation>
	</template>

</service>

</service_bundle>
