# usr/src/cmd/initpkg/inittab
# CDDL HEADER START
#
# The contents of this file are subject to the terms of the
# Common Development and Distribution License (the "License").
# You may not use this file except in compliance with the License.
#
# You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
# or http://www.opensolaris.org/os/licensing.
# See the License for the specific language governing permissions
# and limitations under the License.
#
# When distributing Covered Code, include this CDDL HEADER in each
# file and include the License file at usr/src/OPENSOLARIS.LICENSE.
# If applicable, add the following below this CDDL HEADER, with the
# fields enclosed by brackets "[]" replaced with your own identifying
# information: Portions Copyright [yyyy] [name of copyright owner]
#
# CDDL HEADER END
#
# Copyright (c) 1988, 2010, Oracle and/or its affiliates. All rights reserved.
#
# The /etc/inittab file controls the configuration of init(8); for more
# information refer to init(8) and inittab(5).  It is no longer
# necessary to edit inittab(5) directly; administrators should use the
# Solaris Service Management Facility (SMF) to define services instead.
# Refer to smf(7) and the System Administration Guide for more
# information on SMF.
#
# For modifying parameters passed to ttymon, use svccfg(8) to modify
# the SMF repository. For example:
#
#	# svccfg
#	svc:> select system/console-login
#	svc:/system/console-login> setprop ttymon/terminal_type = "xterm"
#	svc:/system/console-login> exit
#
ap::sysinit:/sbin/autopush -f /etc/iu.ap
sp::sysinit:/sbin/soconfig -d /etc/sock2path.d
smf::sysinit:/lib/svc/bin/svc.startd	>/dev/msglog 2<>/dev/msglog </dev/console
p3:s1234:powerfail:/usr/sbin/shutdown -y -i5 -g0 >/dev/msglog 2<>/dev/msglog
