# usr/src/cmd/fm/fmd/common/fmd.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
 Copyright (c) 2004, 2010, Oracle and/or its affiliates. All rights reserved.

 CDDL HEADER START

 The contents of this file are subject to the terms of the
 Common Development and Distribution License (the "License").
 You may not use this file except in compliance with the License.

 You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
 or http://www.opensolaris.org/os/licensing.
 See the License for the specific language governing permissions
 and limitations under the License.

 When distributing Covered Code, include this CDDL HEADER in each
 file and include the License file at usr/src/OPENSOLARIS.LICENSE.
 If applicable, add the following below this CDDL HEADER, with the
 fields enclosed by brackets "[]" replaced with your own identifying
 information: Portions Copyright [yyyy] [name of copyright owner]

 CDDL HEADER END

	NOTE:  This service manifest is not editable; its contents will
	be overwritten by package or patch operations, including
	operating system upgrade.  Make customizations in a different
	file.
-->

<service_bundle type='manifest' name='SUNWckr:fmd'>

<service
	name='system/fmd'
	type='service'
	version='1'>

	<create_default_instance enabled='true' />

	<single_instance />

	<!--
	  Fault Manager Dependencies, Part 1

	    At present, the fmd.xml manifest is delivered outside of SUNWfmd
	    because it is a usr package.  Therefore we insert a dependency on
	    the fmd binary being present (i.e. SUNWfmd has been installed).
	-->
	<dependency
	    name='SUNWfmd'
	    grouping='require_all'
	    restart_on='none'
	    type='path'>
		<service_fmri value='file://localhost/usr/lib/fm/fmd/fmd' />
	</dependency>

	<!--
	  Fault Manager Dependencies, Part 2

	    fmd should be started after /var is mounted and rpcbind is
	    running.  It is resilient to the the failure and/or
	    interruption of rpcbind, and therefore does not restart under any
	    circumstance.

	    The dumpadm and sysevent services (global zone only) also list the
	    fmd service as a dependent, with restart_on='none' for both.
	-->
	<dependency
	    name='startup_req'
	    grouping='require_all'
	    restart_on='none'
	    type='service'>
		<service_fmri value='svc:/system/filesystem/minimal' />
	</dependency>

	<dependency
	    name='startup_opt'
	    grouping='optional_all'
	    restart_on='none'
	    type='service'>
		<service_fmri value='svc:/network/rpc/bind' />
	</dependency>

	<exec_method
	    type='method'
	    name='start'
	    exec='/usr/lib/fm/fmd/fmd'
	    timeout_seconds='0' />

	<exec_method
	    type='method'
	    name='stop'
	    exec=':kill'
	    timeout_seconds='60' />

	<stability value='Evolving' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
				Fault Manager
			</loctext>
		</common_name>

		<documentation>
			<manpage
			    title='fmd'
			    section='8'
			    manpath='/usr/share/man' />
		</documentation>
	</template>
</service>

</service_bundle>
