# usr/src/cmd/xntpd/xntpd/ntp.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
    Copyright 2009 Sun Microsystems, Inc.  All rights reserved.
    Use is subject to license terms.

    NOTE:  This service manifest is not editable; its contents will
    be overwritten by package or patch operations, including
    operating system upgrade.  Make customizations in a different
    file.
-->

<service_bundle type='manifest' name='SUNWntpr:xntpd'>

<service
	name='network/ntp'
	type='service'
	version='1'>

	<create_default_instance enabled='false' />

	<single_instance />

	<!--
		In its multicast mode, xntpd(1M) requires the default
		multicast route to be established.
	-->
	<dependency name='network'
	    grouping='require_any'
	    restart_on='error'
	    type='service'>
		<service_fmri value='svc:/network/service' />
	</dependency>

	<dependent
		name='ntp_multi-user'
		grouping='optional_all'
		restart_on='none'>
		<service_fmri value='svc:/milestone/multi-user' />
	</dependent>

	<!--
		The timeout needs to be large enough to handle contacting
		the peer system.
	-->
	<exec_method
	    type='method'
	    name='start'
	    exec='/lib/svc/method/xntp'
	    timeout_seconds='1800' />

	<exec_method
	    type='method'
	    name='stop'
	    exec=':kill'
	    timeout_seconds='60' />

	<property_group name='firewall_context' type='com.sun,fw_definition'>
		<propval name='name' type='astring' value='ntp' />
		<propval name='ipf_method' type='astring'
			value='/lib/svc/method/xntp ipfilter' />
	</property_group>

	<property_group name='firewall_config' type='com.sun,fw_configuration'>
		<propval name='policy' type='astring' value='use_global' />
		<propval name='apply_to' type='astring' value='' />
		<propval name='exceptions' type='astring' value='' />
		<propval name='value_authorization' type='astring'
			value='solaris.smf.value.firewall.config' />
	</property_group>

	<stability value='Unstable' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
			Network Time Protocol (NTP)
			</loctext>
		</common_name>
		<documentation>
			<manpage title='xntpd' section='1M' />
			<manpage title='ntpdate' section='1M' />
			<manpage title='ntpq' section='1M' />
		</documentation>
	</template>
</service>

</service_bundle>
