# usr/src/cmd/cmd-inet/usr.sbin/in.routed/route.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
 Copyright 2009 Sun Microsystems, Inc.  All rights reserved.
 Use is subject to license terms.

 Copyright 2016 Hans Rosenfeld <rosenfeld@grumpf.hope-2000.org>

 CDDL HEADER START

 The contents of this file are subject to the terms of the
 Common Development and Distribution License (the "License").
 You may not use this file except in compliance with the License.

 You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
 or http://www.opensolaris.org/os/licensing.
 See the License for the specific language governing permissions
 and limitations under the License.

 When distributing Covered Code, include this CDDL HEADER in each
 file and include the License file at usr/src/OPENSOLARIS.LICENSE.
 If applicable, add the following below this CDDL HEADER, with the
 fields enclosed by brackets "[]" replaced with your own identifying
 information: Portions Copyright [yyyy] [name of copyright owner]

 CDDL HEADER END

	NOTE:  This service manifest is not editable; its contents will
	be overwritten by package or patch operations, including
	operating system upgrade.  Make customizations in a different
	file.
-->

<service_bundle type='manifest' name='SUNWroutr:route'>

<service
	name='network/routing/route'
	type='service'
	version='1'>

	<single_instance />

	<instance name='default' enabled='false' >

	<!-- in.routed will not run unless routing-setup has run -->
	<dependency
		name='network_routing_setup'
		grouping='require_all'
		restart_on='refresh'
		type='service'>
		<service_fmri value='svc:/network/routing-setup' />
	</dependency>

	<exec_method
		type='method'
		name='start'
		exec='/lib/svc/method/svc-route'
		timeout_seconds='60'>
		<method_context>
			<method_credential user='root' group='root'
privileges='basic,proc_owner,proc_fork,proc_exec,proc_info,proc_session,file_chown,sys_ip_config,net_privaddr,net_icmpaccess,net_rawaccess'/>
		</method_context>
	</exec_method>

	<exec_method
		type='method'
		name='stop'
		exec=':kill'
		timeout_seconds='60'>
		<method_context>
			<method_credential user='root' group='root'/>
		</method_context>
	</exec_method>

	<!-- to start stop routing services -->
	<property_group name='general' type='framework'>
		<propval name='action_authorization' type='astring'
			value='solaris.smf.manage.routing' />
		<propval name='value_authorization' type='astring'
			value='solaris.smf.manage.routing' />
	</property_group>

	<!-- Properties in this group are used by routeadm (1M) -->
	<property_group name='routeadm' type='application'>
		<stability value='Unstable' />
		<propval name='protocol' type='astring' value='ipv4' />
		<propval name='daemon' type='astring'
			value='/usr/sbin/in.routed' />
		<propval name='value_authorization' type='astring'
			value='solaris.smf.value.routing' />
	</property_group>

	<!-- Properties in this group are modifiable via routeadm (1M) -->
	<property_group name='routing' type='application'>
		<stability value='Evolving' />

		<!-- Equivalent to -A option if false -->
		<propval name='ignore_auth' type='boolean' value='true' />

		<!-- Equivalent to -F specificied route minimization list -->
		<propval name='minimize_routes' type='astring' value=''/>

		<!-- Equivalent to -g option if true -->
		<propval name='offer_default_route' type='boolean'
		    value='false' />

		<!-- Equivalent to -h option if false -->
		<propval name='advertise_host_routes' type='boolean'
		    value='true' />

		<!-- Equivalent to -m option if true -->
		<propval name='advertise_host_routes_primary' type='boolean'
		    value='false' />

		<!-- Equivalent to -n option if false -->
		<propval name='install_routes' type='boolean' value='true' />

		<!-- Equivalent to -P specified parameter list -->
		<propval name='parameters' type='astring' value='' />

		<!-- Equivalent to the -q option if true -->
		<propval name='quiet_mode' type='boolean' value='false' />

		<!-- Equivalent to -s option if true -->
		<propval name='supply_routes' type='boolean' value='false' />	

		<!-- Equivalent to -S option if true -->
		<propval name='default_routes_only' type='boolean'
			value='false' />

		<!-- Equivalent to -T tracefile option -->
		<propval name='log_file' type='astring' value='' />

		<!-- Equivalent to -v/-z option if true -->
		<propval name='debug' type='boolean' value='false' />

		<propval name='value_authorization' type='astring'
			value='solaris.smf.value.routing' />
	</property_group>

	<property_group name='firewall_context' type='com.sun,fw_definition'>
		<propval name='name' type='astring' value='route' />
		<propval name='ipf_method' type='astring'
			value='/lib/svc/method/svc-route ipfilter' />
	</property_group>

	<property_group name='firewall_config' type='com.sun,fw_configuration'>
		<propval name='policy' type='astring' value='use_global' />
		<propval name='block_policy' type='astring'
			value='use_global' />
		<propval name='apply_to' type='astring' value='' />
		<propval name='exceptions' type='astring' value='' />
		<propval name='target' type='astring' value='' />
		<propval name='value_authorization' type='astring'
			value='solaris.smf.value.firewall.config' />
	</property_group>

	<template>
		<common_name>
			<loctext xml:lang='C'>
				in.routed network routing daemon
			</loctext>
		</common_name>
		<documentation>
			<manpage title='in.routed' section='1M'
				manpath='/usr/share/man' />
		</documentation>
	</template>
	</instance>
	<stability value='Unstable' />
</service>

</service_bundle>
