# usr/src/lib/print/mod_ipp/ipp-listener.xml
<?xml version="1.0"?>
<!--
CDDL HEADER START

The contents of this file are subject to the terms of the
Common Development and Distribution License (the "License").
You may not use this file except in compliance with the License.

You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
or http://www.opensolaris.org/os/licensing.
See the License for the specific language governing permissions
and limitations under the License.

When distributing Covered Code, include this CDDL HEADER in each
file and include the License file at usr/src/OPENSOLARIS.LICENSE.
If applicable, add the following below this CDDL HEADER, with the
fields enclosed by brackets "[]" replaced with your own identifying
information: Portions Copyright [yyyy] [name of copyright owner]

CDDL HEADER END
-->
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
 Copyright 2009 Sun Microsystems, Inc.  All rights reserved.
 Use is subject to license terms.
-->

<service_bundle type='manifest' name='SUNWipplr:ipp-listener'>

<service
    name='application/print/ipp-listener'
    type='service'
    version='1'>

	<create_default_instance enabled='false' />

	<dependency name='print-service'
	    grouping='require_any'
	    restart_on='refresh'
	    type='service'>
	    <service_fmri value='svc:/application/print/server' />
	</dependency>

	<!-- Change /usr/apache2/... if Apache 2.x is elsewhere. -->
	<exec_method
	    type='method'
	    name='start'
	    exec='/usr/apache2/2.2/bin/httpd -f /etc/lp/ipp/httpd-standalone-ipp.conf'
	    timeout_seconds='10' />

	<exec_method
	    type='method'
	    name='stop'
	    exec='/bin/pkill -f httpd-standalone-ipp.conf'
	    timeout_seconds='5' />

	<property_group name='general' type='framework'>
		<!-- to start/stop IPP listening service-->
		<propval name='action_authorization' type='astring'
			value='solaris.print.admin' />
		<propval name='value_authorization' type='astring'
			value='solaris.print.admin' />
	</property_group>

	<property_group name='firewall_context' type='com.sun,fw_definition'>
		<propval name='ipf_method' type='astring'
			value='/lib/svc/method/print-svc ipfilter svc:/application/print/server:default' />
	</property_group>

	<stability value='Unstable' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
				Internet Print Protocol Listening Service
			</loctext>
		</common_name>
		<documentation>
			<manpage title='mod_ipp' section='4'
				manpath='/usr/share/man' />
		</documentation>
	</template>
</service>

</service_bundle>
