# usr/src/cmd/dumpadm/dumpadm.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
 Copyright (c) 2004, 2010, Oracle and/or its affiliates. All rights reserved.

 CDDL HEADER START

 The contents of this file are subject to the terms of the
 Common Development and Distribution License (the "License").
 You may not use this file except in compliance with the License.

 You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
 or http://www.opensolaris.org/os/licensing.
 See the License for the specific language governing permissions
 and limitations under the License.

 When distributing Covered Code, include this CDDL HEADER in each
 file and include the License file at usr/src/OPENSOLARIS.LICENSE.
 If applicable, add the following below this CDDL HEADER, with the
 fields enclosed by brackets "[]" replaced with your own identifying
 information: Portions Copyright [yyyy] [name of copyright owner]

 CDDL HEADER END

	NOTE:  This service manifest is not editable; its contents will
	be overwritten by package or patch operations, including
	operating system upgrade.  Make customizations in a different
	file.

	Service manifest for dumpadm.
-->

<service_bundle type='manifest' name='SUNWckr:dumpadm'>

<service
	name='system/dumpadm'
	type='service'
	version='1'>

	<create_default_instance enabled='true' />

	<single_instance />

	<!--
		The dumpadm(1M) service reconfigures the kernel's dump settings
		from /etc/dumpadm.conf and then invokes savecore(1M) to locate
		and retrieve any saved crash dumps.  It therefore needs to have
		local filesystems mounted, and must also have an *optional*
		dependency on nfs.client and autofs in case dumpadm -s has been
		used to set the savecore directory to an NFS or AutoFS mount.
	-->

	<dependency
		name='local-filesystems'
		type='service'
		grouping='require_all'
		restart_on='none'>
		<service_fmri value='svc:/system/filesystem/local' />
	</dependency>

	<dependency
		name='remote-filesystems'
		type='service'
		grouping='optional_all'
		restart_on='none'>
		<service_fmri value='svc:/network/nfs/client' />
		<service_fmri value='svc:/system/filesystem/autofs' />
	</dependency>

	<!--
		The fmd service may exist in a zone, but only in the
		global zone (which dumpadm is defined in) should it be
		dependent on dumpadm configuration.
	-->
	<dependent
		name='dumpadm-fmd'
		grouping='require_all'
		restart_on='none'>
		<service_fmri value='svc:/system/fmd' />
	</dependent>

	<exec_method
		type='method'
		name='start'
		exec='/lib/svc/method/svc-dumpadm'
		timeout_seconds='0' />

	<exec_method
		type='method'
		name='stop'
		exec=':true'
		timeout_seconds='0' />

	<property_group name='startd' type='framework'>
		<propval name='duration' type='astring' value='transient' />
	</property_group>

	<stability value='Unstable' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
				system crash dump configuration
			</loctext>
		</common_name>

		<documentation>
			<manpage
				title='dumpadm'
				section='1M'
				manpath='/usr/share/man' />
		</documentation>
	</template>
</service>

</service_bundle>
