# usr/src/cmd/fs.d/nfs/svc/status.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
 Copyright 2009 Sun Microsystems, Inc.  All rights reserved.
 Use is subject to license terms.

 CDDL HEADER START

 The contents of this file are subject to the terms of the
 Common Development and Distribution License (the "License").
 You may not use this file except in compliance with the License.

 You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
 or http://www.opensolaris.org/os/licensing.
 See the License for the specific language governing permissions
 and limitations under the License.

 When distributing Covered Code, include this CDDL HEADER in each
 file and include the License file at usr/src/OPENSOLARIS.LICENSE.
 If applicable, add the following below this CDDL HEADER, with the
 fields enclosed by brackets "[]" replaced with your own identifying
 information: Portions Copyright [yyyy] [name of copyright owner]

 CDDL HEADER END

	NOTE:  This service manifest is not editable; its contents will
	be overwritten by package or patch operations, including
	operating system upgrade.  Make customizations in a different
	file.

	Note: if this service is modified to consist of anything other
	than a single instance named 'default', you must make changes to
	$SRC/head/rpcsvc/daemon_utils.h and libnsl:open_daemon_lock().
-->

<!-- Copyright (c) 2016 by Delphix. All rights reserved. -->

<service_bundle type='manifest' name='SUNWnfscr:nfs-status'>

<service
	name='network/nfs/status'
	type='service'
	version='1'>

	<dependency name='network'
	    grouping='require_any'
	    restart_on='error'
	    type='service'>
		<service_fmri value='svc:/milestone/network' />
	</dependency>

	<dependency name='rpcbind'
	    grouping='require_all'
	    restart_on='restart'
	    type='service'>
		<service_fmri value='svc:/network/rpc/bind' />
	</dependency>

	<dependency name='filesystem-local'
	    grouping='require_all'
	    restart_on='error'
	    type='service'>
	    	<service_fmri value='svc:/system/filesystem/local' />
	</dependency>

	<exec_method
	    type='method'
	    name='start'
	    exec='/usr/lib/nfs/statd'
	    timeout_seconds='60' />

	<exec_method
	    type='method'
	    name='stop'
	    exec=':kill'
	    timeout_seconds='60' />

	<property_group name='firewall_context' type='com.sun,fw_definition'>
		<propval name='isrpc' type='boolean' value='true' />
		<propval name='name' type='astring' value='status' />
		<propval name='ipf_method' type='astring'
			value='/lib/svc/method/nfs-server ipfilter' />
	</property_group>

	<property_group name='application' type='framework'>
		<stability value='Evolving' />
		<propval name='auto_enable' type='boolean' value='true' />
	</property_group>

	<instance name='default' enabled='false'>
		<property_group name='nfs-props' type='com.oracle.nfs,props'>
			<propval name='statd_port' type='integer' value='0'/>
		</property_group>
	</instance>

	<stability value='Stable' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
				NFS status monitor
			</loctext>
		</common_name>
		<documentation>
			<manpage title='statd' section='1M'
				manpath='/usr/share/man' />
		</documentation>
	</template>
</service>

</service_bundle>
