# usr/src/cmd/ipf/svc/ipfilter.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
 Copyright 2009 Sun Microsystems, Inc.  All rights reserved.
 Use is subject to license terms.

 CDDL HEADER START

 The contents of this file are subject to the terms of the
 Common Development and Distribution License (the "License").
 You may not use this file except in compliance with the License.

 You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
 or http://www.opensolaris.org/os/licensing.
 See the License for the specific language governing permissions
 and limitations under the License.

 When distributing Covered Code, include this CDDL HEADER in each
 file and include the License file at usr/src/OPENSOLARIS.LICENSE.
 If applicable, add the following below this CDDL HEADER, with the
 fields enclosed by brackets "[]" replaced with your own identifying
 information: Portions Copyright [yyyy] [name of copyright owner]

 CDDL HEADER END

	NOTE:  This service description is not editable; its contents
	may be overwritten by package or patch operations, including
	operating system upgrade.  Make customizations in a different
	file.

	Service manifest for the ipfilter service.
-->

<service_bundle type='manifest' name='SUNWipfr:ipfilter'>

<service
	name='network/ipfilter'
	type='service'
	version='1'>

	<create_default_instance enabled='false' />

	<single_instance />

	<dependency 
	    name='config'
	    grouping='require_all'
	    restart_on='restart'
	    type='path'>
		<service_fmri value='file://localhost/etc/ipf/ipf.conf' />
	</dependency>

	<dependency
	    name='filesystem'
	    grouping='require_all'
	    restart_on='none'
	    type='service'>
		<service_fmri value='svc:/system/filesystem/minimal' />
	</dependency>

	<dependency
	    name='physical'
	    grouping='require_all'
	    restart_on='restart'
	    type='service'>
		<service_fmri value='svc:/network/physical' />
	</dependency>

	<dependency
	    name='identity'
	    grouping='require_all'
	    restart_on='restart'
	    type='service'>
		<service_fmri value='svc:/system/identity:node' />
	</dependency>

	<dependent
	    name='network'
	    grouping='optional_all'
	    restart_on='restart'>
		<service_fmri value='svc:/milestone/network' />
	</dependent>

	<exec_method
		type='method'
		name='stop'
		exec='/lib/svc/method/ipfilter stop'
		timeout_seconds='60' >
	</exec_method>

	<exec_method
		type='method'
		name='start'
		exec='/lib/svc/method/ipfilter start'
		timeout_seconds='30' >
	</exec_method>

	<exec_method
		type='method'
		name='refresh'
		exec='/lib/svc/method/ipfilter reload'
		timeout_seconds='30' >
	</exec_method>

	<property_group
		name='startd'
		type='framework'>
		<propval name='duration' type='astring' value='transient' />
	</property_group>

	<stability value='Unstable' />

	<template>
		<common_name>
			<loctext xml:lang='C'>IP Filter</loctext>
		</common_name>
		<description>
			<loctext xml:lang='C'>
			Solaris IP Filter - host-based firewall
			</loctext>
		</description>
		<documentation>
			<manpage title='ipfilter' section='5'
				manpath='/usr/share/man' />
		</documentation>
	</template>
</service>

</service_bundle>
