# usr/src/cmd/iscsi/iscsitgtd/iscsi_target.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
 Copyright 2007 Sun Microsystems, Inc.  All rights reserved.
 Use is subject to license terms.

 CDDL HEADER START

  The contents of this file are subject to the terms of the
  Common Development and Distribution License (the "License").
  You may not use this file except in compliance with the License.
 
  You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
  or http://www.opensolaris.org/os/licensing.
  See the License for the specific language governing permissions
  and limitations under the License.
 
  When distributing Covered Code, include this CDDL HEADER in each
  file and include the License file at usr/src/OPENSOLARIS.LICENSE.
  If applicable, add the following below this CDDL HEADER, with the
  fields enclosed by brackets "[]" replaced with your own identifying
  information: Portions Copyright [yyyy] [name of copyright owner]

 CDDL HEADER END

	ident	"%Z%%M%	%I%	%E% SMI"


	NOTE:  This service manifest is not editable; its contents will
	be overwritten by package or patch operations, including
	operating system upgrade.  Make customizations in a different
	file.
-->

<service_bundle type='manifest' name='SUNWiscsitgtu:iscsitgt'>

<service
	name='system/iscsitgt'
	type='service'
	version='1'>

	<create_default_instance enabled='false' />

	<single_instance/>

	<!-- We need name resolution and full filesystem access -->
	<dependency
		name='milestone'
		grouping='require_all'
		restart_on='none'
		type='service'>
		<service_fmri value='svc:/milestone/multi-user' />
	</dependency>


	<exec_method
		type='method'
		name='start'
		exec='/lib/svc/method/svc-iscsitgt %m'
		timeout_seconds='60'>
		<method_context>
			<method_credential
				user='root'
				group='sys'
				privileges='basic,sys_config,net_rawaccess,sys_mount,file_dac_write,sys_devices,proc_setid' />
		</method_context>
	</exec_method>

	<exec_method
		type='method'
		name='stop'
		exec='/lib/svc/method/svc-iscsitgt %m %{restarter/contract}'
		timeout_seconds='60'>
		<method_context>
			<method_credential
				user='root'
				group='sys'
				privileges='basic,sys_config,net_rawaccess,sys_mount,file_dac_write,sys_devices,proc_setid' />
		</method_context>
	</exec_method>

	<property_group name='general' type='framework'>
		<propval name='action_authorization' type='astring'
			value='solaris.smf.manage.iscsitgt' />
		<!-- To read and modify protected properties -->
		<propval name='read_authorization' type='astring'
			value='solaris.smf.read.iscsitgt' />
		<propval name='value_authorization' type='astring'
			value='solaris.smf.value.iscsitgt' />
		<propval name='modify_authorization' type='astring'
			value='solaris.smf.modify.iscsitgt' />
	</property_group>

	<property_group name='application' type='framework'>
		<stability value='Unstable' />
		<propval name='auto_enable' type='boolean' value='true' />
	</property_group>

	<stability value='Unstable' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
				iSCSI Target
			</loctext>
		</common_name>
		<documentation>
			<manpage title='iscsitgtd' section='1M'
				manpath='/usr/share/man' />
			<manpage title='iscsitadm' section='1M'
				manpath='/usr/share/man' />
		</documentation>
	</template>
</service>

</service_bundle>
