# deleted_files/usr/src/cmd/ipf/svc/pfil.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
	Copyright 2004 Sun Microsystems, Inc.  All rights reserved.
	Use is subject to license terms.

	ident	"%Z%%M%	%I%	%E% SMI"

	NOTE:  This service manifest is not editable; its contents will
	be overwritten by package or patch operations, including
	operating system upgrade.  Make customizations in a different
	file.

	Service manifest for the pfil (packet filter) service.
-->

<service_bundle type='manifest' name='SUNWipfr:pfil'>

<service
	name='network/pfil'
	type='service'
	version='1'>

	<create_default_instance enabled='true' />

	<single_instance />

	<dependent
	    name='pfil_network'
	    grouping='optional_all'
	    restart_on='none'>
		<service_fmri value='svc:/network/physical' />
	</dependent>

	<dependent
	    name='pfil_sysid'
	    grouping='optional_all'
	    restart_on='none'>
		<service_fmri value='svc:/system/sysidtool:net' />
	</dependent>

	<!--
		The stop method really should deconfigure pfil from sad(7D),
		but autopush(1M) doesn't have the ability to remove entries
		based on a file in the same format as it uses to add them.
	-->
	<exec_method
		type='method'
		name='stop'
		exec=':true'
		timeout_seconds='0' >
	</exec_method>

	<exec_method
		type='method'
		name='start'
		exec='/lib/svc/method/pfil start'
		timeout_seconds='0' >
	</exec_method>

	<property_group
		name='startd'
		type='framework'>
		<propval name='duration' type='astring' value='transient' />
	</property_group>

	<stability value='Unstable' />

	<template>
		<common_name>
			<loctext xml:lang='C'>packet filter</loctext>
		</common_name>
		<description>
			<loctext xml:lang='C'>
			Packet filter interface; autopushes the pfil STREAMS
			module on network devices.
			</loctext>
		</description>
		<documentation>
			<manpage title='ipfilter' section='5'
				manpath='/usr/share/man' />
		</documentation>
	</template>
</service>

</service_bundle>
