# usr/src/cmd/svc/milestone/identity.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
 Copyright 2005 Sun Microsystems, Inc.  All rights reserved.
 Use is subject to license terms.

 CDDL HEADER START

 The contents of this file are subject to the terms of the
 Common Development and Distribution License, Version 1.0 only
 (the "License").  You may not use this file except in compliance
 with the License.

 You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
 or http://www.opensolaris.org/os/licensing.
 See the License for the specific language governing permissions
 and limitations under the License.

 When distributing Covered Code, include this CDDL HEADER in each
 file and include the License file at usr/src/OPENSOLARIS.LICENSE.
 If applicable, add the following below this CDDL HEADER, with the
 fields enclosed by brackets "[]" replaced with your own identifying
 information: Portions Copyright [yyyy] [name of copyright owner]

 CDDL HEADER END

	NOTE:  This service manifest is not editable; its contents will
	be overwritten by package or patch operations, including
	operating system upgrade.  Make customizations in a different
	file.
-->

<service_bundle type='manifest' name='SUNWcsr:identity'>

<service
	name='system/identity'
	type='service'
	version='1'>

	<dependency
		name='loopback-network'
		grouping='require_any'
		restart_on='none'
		type='service'>
		<service_fmri value='svc:/network/loopback' />
	</dependency>

	<dependency
		name='physical-network'
		grouping='optional_all'
		restart_on='none'
		type='service'>
		<service_fmri value='svc:/network/physical' />
	</dependency>

	<exec_method
		type='method'
		name='stop'
		exec=':true'
		timeout_seconds='60' />

	<property_group name='startd' type='framework'>
		<propval name='duration' type='astring' value='transient' />
	</property_group>

	<instance name='node' enabled='true'>

		<exec_method
			type='method'
			name='start'
			exec='/lib/svc/method/identity-node'
			timeout_seconds='60' />

		<template>
			<common_name>
				<loctext xml:lang='C'>
					system identity (nodename)
				</loctext>
			</common_name>
			<documentation>
				<manpage title='nodename' section='5'
					manpath='/usr/share/man' />
			</documentation>
		</template>

	</instance>

	<instance name='domain' enabled='false'>

		<dependency
			name='fs'
			grouping='require_all'
			restart_on='none'
			type='service'>
			<service_fmri value='svc:/system/filesystem/minimal' />
		</dependency>

		<exec_method
			type='method'
			name='start'
			exec='/lib/svc/method/identity-domain'
			timeout_seconds='60' />

		<template>
			<common_name>
				<loctext xml:lang='C'>
					system identity (domainname)
				</loctext>
			</common_name>
			<documentation>
				<manpage title='domainname' section='8'
					manpath='/usr/share/man' />
				<manpage title='defaultdomain' section='5'
					manpath='/usr/share/man' />
			</documentation>
		</template>

	</instance>

	<stability value='Unstable' />

</service>

</service_bundle>
