# usr/src/cmd/svc/milestone/network-service.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
 Copyright 2009 Sun Microsystems, Inc.  All rights reserved.
 Use is subject to license terms.

 CDDL HEADER START

 The contents of this file are subject to the terms of the
 Common Development and Distribution License (the "License").
 You may not use this file except in compliance with the License.

 You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
 or http://www.opensolaris.org/os/licensing.
 See the License for the specific language governing permissions
 and limitations under the License.

 When distributing Covered Code, include this CDDL HEADER in each
 file and include the License file at usr/src/OPENSOLARIS.LICENSE.
 If applicable, add the following below this CDDL HEADER, with the
 fields enclosed by brackets "[]" replaced with your own identifying
 information: Portions Copyright [yyyy] [name of copyright owner]

 CDDL HEADER END

	NOTE:  This service manifest is not editable; its contents will
	be overwritten by package or patch operations, including
	operating system upgrade.  Make customizations in a different
	file.
-->

<service_bundle type='manifest' name='SUNWcsr:network-service'>

<!--
  network/service is the fourth service containing aggregated TCP/IP
  service initialization.  It will decompose into its constituent
  services over time.
-->

<service
	name='network/service'
	type='service'
	version='1'>

	<create_default_instance enabled='true' />

	<dependency
		name='init'
		grouping='optional_all'
		restart_on='none'
		type='service'>
		<service_fmri value='svc:/network/initial' />
	</dependency>

	<!--
	  We require the root, usr and minimal filesystems, since
	  this service writes to/reads from /tmp, /etc, /var and
	  runs programs under /usr.
	-->
	<dependency
		name='filesystem'
		grouping='require_all'
		restart_on='none'
		type='service'>
		<service_fmri value='svc:/system/filesystem/root' />
		<service_fmri value='svc:/system/filesystem/usr' />
		<service_fmri value='svc:/system/filesystem/minimal' />
	</dependency>

	<dependency
		name='nis_server'
		grouping='optional_all'
		restart_on='none'
		type='service'>
		<service_fmri value='svc:/network/nis/server' />
	</dependency>

	<dependency
		name='nis_client'
		grouping='optional_all'
		restart_on='none'
		type='service'>
		<service_fmri value='svc:/network/nis/client' />
	</dependency>

	<!--
	  DNS is potentially configured by the DHCP actions in
	  network/service, and cannot presently be used to store data
	  used by this service.
	-->

	<exec_method
		type='method'
		name='start'
		exec='/lib/svc/method/net-svc start'
		timeout_seconds='600' />

	<exec_method
		type='method'
		name='stop'
		exec=':true'
		timeout_seconds='3' />

	<property_group name='startd' type='framework'>
		<propval name='duration' type='astring' value='transient' />
	</property_group>

	<stability value='Unstable' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
				layered network services
			</loctext>
		</common_name>
		<description>
			<loctext xml:lang='C'>
				Network infrastructure services
				requiring name service availability.
			</loctext>
		</description>
		<documentation>
			<manpage title='ifconfig' section='8'
				manpath='/usr/share/man' />
		</documentation>
	</template>
</service>

</service_bundle>
