# usr/src/cmd/svc/milestone/network-early-admin.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
 Copyright 2019, Joyent, Inc.

  This file and its contents are supplied under the terms of the
  Common Development and Distribution License ("CDDL"), version 1.0.
  You may only use this file in accordance with the terms of version
  1.0 of the CDDL.

  A full copy of the text of the CDDL should have accompanied this
  source.  A copy of the CDDL is also available via the Internet at
  http://www.illumos.org/license/CDDL.

-->

<service_bundle type='manifest' name='network-early-admin'>

<service
	name='network/early-admin'
	type='service'
	version='1'>

	<!-- ifconfig needs loopback for IPC with dhcpagent -->
	<dependency
		name='network'
		grouping='require_all'
		restart_on='none'
		type='service'>
		<service_fmri value='svc:/network/loopback' />
		<service_fmri value='svc:/network/datalink-management' />
		<service_fmri value='svc:/network/ip-interface-management' />
	</dependency>

	<instance name='default' enabled='true'>

	<exec_method
		type='method'
		name='start'
		exec='/lib/svc/method/net-early-admin'
		timeout_seconds='600' />

	<exec_method
		type='method'
		name='stop'
		exec=':true'
		timeout_seconds='3' />

	<property_group name='startd' type='framework'>
		<propval name='duration' type='astring' value='transient' />
	</property_group>

	<template>
		<common_name>
			<loctext xml:lang='C'>
				Triton admin network on compute nodes
			</loctext>
		</common_name>
		<documentation>
			<manpage title='ifconfig' section='1M'
				manpath='/usr/share/man' />
			<manpage title='dladm' section='1M'
				manpath='/usr/share/man' />
		</documentation>
	</template>

	</instance>

	<stability value='Unstable' />

</service>

</service_bundle>
