# usr/src/cmd/cmd-inet/usr.sbin/in.ftpd/ftp.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>

<!--
	Copyright 2004 Sun Microsystems, Inc.  All rights reserved.
	Use is subject to license terms.

	ident	"%Z%%M%	%I%	%E% SMI"

	NOTE:  This service manifest is not editable; its contents will
	be overwritten by package or patch operations, including
	operating system upgrade.  Make customizations in a different
	file.

	Service manifest for in.ftpd
-->

<service_bundle type='manifest' name='SUNWftpr:ftp'>

<service
	name='network/ftp'
	type='service'
	version='1'>

	<create_default_instance enabled='false' />

	<restarter>
		<service_fmri value='svc:/network/inetd:default' />
	</restarter>

	<exec_method
		type='method'
		name='inetd_start'
		exec='/usr/sbin/in.ftpd -a'
		timeout_seconds='0'>
		<method_context>
			<method_credential user='root' group='root' />
		</method_context>
	</exec_method>

	<exec_method
		type='method'
		name='inetd_disable'
		exec=':kill'
		timeout_seconds='0'>
	</exec_method>

	<property_group name='inetd' type='framework'>
		<stability value='Evolving' />
		<propval name='endpoint_type' type='astring' value='stream' />
		<propval name='name' type='astring' value='ftp' />
		<propval name='proto' type='astring' value='tcp6' />
		<propval name='wait' type='boolean' value='false' />
		<propval name='isrpc' type='boolean' value='false' />
	</property_group>

	<stability value='Unstable' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
				FTP server
			</loctext>
		</common_name>
		<documentation>
			<manpage title='in.ftpd' section='1M'
			    manpath='/usr/share/man' />
			<manpage title='ftpd' section='1M'
			    manpath='/usr/share/man' />
		</documentation>
	</template>

</service>

</service_bundle>
