# usr/src/cmd/varpd/varpd.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
This file and its contents are supplied under the terms of the
Common Development and Distribution License ("CDDL"), version 1.0.
You may only use this file in accordance with the terms of version
1.0 of the CDDL.

A full copy of the text of the CDDL should have accompanied this
source.  A copy of the CDDL is also available via the Internet at
http://www.illumos.org/license/CDDL.

Copyright 2018, Joyent, Inc.
-->

<service_bundle type="manifest" name="illumos:varpd" >

	<service name="network/varpd" type="service" version="1" >

		<create_default_instance enabled="true" />

		<single_instance/>

		<dependency name="varpd-network-physical"
			grouping="require_all"
			restart_on="none"
			type="service">
			<service_fmri value="svc:/network/physical:default" />
		</dependency>

		<dependency name="varpd-device-local"
			grouping="require_all"
			restart_on="none"
			type="service">
			<service_fmri value="svc:/system/device/local:default" />
		</dependency>

		<exec_method
			type="method"
			name="start"
			exec="/usr/lib/varpd/varpd"
			timeout_seconds="60" />

		<exec_method
			type="method"
			name="stop"
			exec=":kill"
			timeout_seconds="10" />

		<property_group name='varpd' type='application'>
			<property name='include_path' type='astring'>
				<astring_list>
					<value_node value='/usr/lib/varpd'/>
				</astring_list>
			</property>
		</property_group>

		<stability value='Unstable' />

		<template>
			<common_name>
				<loctext xml:lang="C">virtual ARP daemon
				</loctext>
			</common_name>
		</template>
	</service>
</service_bundle>
