# deleted_files/usr/src/cmd/volmgt/etc/volfs.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
	Copyright 2005 Sun Microsystems, Inc.  All rights reserved.
	Use is subject to license terms.

 CDDL HEADER START

 The contents of this file are subject to the terms of the
 Common Development and Distribution License, Version 1.0 only
 (the "License").  You may not use this file except in compliance
 with the License.

 You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
 or http://www.opensolaris.org/os/licensing.
 See the License for the specific language governing permissions
 and limitations under the License.

 When distributing Covered Code, include this CDDL HEADER in each
 file and include the License file at usr/src/OPENSOLARIS.LICENSE.
 If applicable, add the following below this CDDL HEADER, with the
 fields enclosed by brackets "[]" replaced with your own identifying
 information: Portions Copyright [yyyy] [name of copyright owner]

 CDDL HEADER END

	ident	"%Z%%M%	%I%	%E% SMI"

	NOTE:  This service manifest is not editable; its contents will
	be overwritten by package or patch operations, including
	operating system upgrade.  Make customizations in a different
	file.

	Service manifest for volfs.
-->

<service_bundle type='manifest' name='SUNWvolr:volfs'>

<service
	name='system/filesystem/volfs'
	type='service'
	version='1'>

	<create_default_instance enabled='false' />

	<single_instance />

	<dependency name='fs'
		type='service'
		grouping='require_all'
		restart_on='none'>
		<service_fmri value='svc:/system/filesystem/local' />
	</dependency>

	<dependency name='smserver'
		grouping='require_all'
		restart_on='none'
		type='service'>
		<service_fmri value='svc:/network/rpc/smserver' />
	</dependency>

	<dependency
		name='rpcbind'
		type='service'
		grouping='require_all'
		restart_on='restart'>
		<service_fmri value='svc:/network/rpc/bind' />
	</dependency>

	<exec_method
		type='method'
		name='start'
		exec='/lib/svc/method/svc-volfs start'
		timeout_seconds='30'/>

	<!-- If device access problems keep vold from stopping in
	     30 seconds, customize this timeout to a higher value.
	-->
	<exec_method
		type='method'
		name='stop'
		exec=':kill'
		timeout_seconds='30' />

	<exec_method
		type='method'
		name='refresh'
		exec='/lib/svc/method/svc-volfs refresh'
		timeout_seconds='30' />

	<!-- these are passed to vold in the method script -->
	<property_group name='vold' type='application'>
		<propval name='never_writeback_label' type='boolean' value='false' />
		<propval name='root_dir' type='astring' value='' />
		<propval name='config_file' type='astring' value='' />
		<propval name='log_nfs_trace' type='boolean' value='false' />
		<propval name='log_verbose' type='boolean' value='false' />
		<propval name='log_debuglevel' type='count' value='0' />
		<propval name='log_file' type='astring' value='' />
	</property_group>

	<property_group name='general' type='framework'>
		<!-- to start stop vold -->
		<propval name='action_authorization' type='astring'
			value='solaris.smf.manage.volfs' />
		<!-- configure vold -->
		<propval name='value_authorization' type='astring'
			value='solaris.smf.modify.volfs' />
	</property_group>

	<stability value='Unstable' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
			Volume Management filesystem
			</loctext>
		</common_name>
		<documentation>
			<manpage title='volfs' section='7FS' manpath='/usr/man' />
		</documentation>
	</template>

</service>

</service_bundle>
