# usr/src/cmd/cmd-inet/usr.lib/in.ndpd/ndp.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
 Copyright 2007 Sun Microsystems, Inc.  All rights reserved.
 Use is subject to license terms.

 CDDL HEADER START

 The contents of this file are subject to the terms of the
 Common Development and Distribution License (the "License").
 You may not use this file except in compliance with the License.

 You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
 or http://www.opensolaris.org/os/licensing.
 See the License for the specific language governing permissions
 and limitations under the License.

 When distributing Covered Code, include this CDDL HEADER in each
 file and include the License file at usr/src/OPENSOLARIS.LICENSE.
 If applicable, add the following below this CDDL HEADER, with the
 fields enclosed by brackets "[]" replaced with your own identifying
 information: Portions Copyright [yyyy] [name of copyright owner]

 CDDL HEADER END

	NOTE:  This service manifest is not editable; its contents will
	be overwritten by package or patch operations, including
	operating system upgrade.  Make customizations in a different
	file.
-->

<service_bundle type='manifest' name='SUNWroutr:ndp'>

<service
	name='network/routing/ndp'
	type='service'
	version='1'>

	<single_instance />

	<instance name='default' enabled='false' >

	<dependency
		name='network_routing_setup'
		grouping='require_all'
		restart_on='refresh'
		type='service'>
		<service_fmri value='svc:/network/routing-setup' />
	</dependency>

	<exec_method
		type='method'
		name='start'
		exec='/lib/svc/method/svc-ndp'
		timeout_seconds='60'>
		<method_context working_directory='/'>
			<method_credential user='root' group='root' />
		</method_context>
	</exec_method>

	<exec_method
		type='method'
		name='stop'
		exec=':kill'
		timeout_seconds='60'>
		<method_context working_directory='/'>
			<method_credential user='root' group='root'/>
		</method_context>
	</exec_method>

	<!-- to start stop routing services -->
	<property_group name='general' type='framework'>
		<propval name='action_authorization' type='astring'
			value='solaris.smf.manage.routing' />
		<propval name='value_authorization' type='astring'
			value='solaris.smf.manage.routing' />
	</property_group>

	<!-- Properties in this group are used by routeadm(8) -->
	<property_group name='routeadm' type='application'>
		<stability value='Unstable' />
		<!--
			Having the protocol property here ensures that ndp
			appears in the SMF service listing, but does not
			participate in dictating ipv6-routing state (as it
			would if it were set to "ipv6"). If it
			did, ipv6-routing would always be reported as enabled,
			since if any IPv6 interfaces are present, ndp will
			be enabled.
		-->
		<propval name='protocol' type='astring' value='ndp' />
		<propval name='value_authorization' type='astring'
			value='solaris.smf.value.routing' />
	</property_group>


	<!-- Properties in this group are modifiable via routeadm(8) -->
	<property_group name='routing' type='application'>
		<stability value='Evolving' />

		<!-- equivalent to -a option if false -->
		<propval name='stateless_addr_conf' type='boolean'
		    value='true' />

		<!-- equivalent to -d option if true -->
		<propval name='debug' type='boolean' value='false' />

		<!-- equivalent to -f config_file option if set -->
		<propval name='config_file' type='astring' value='' />

		<propval name='value_authorization' type='astring'
			value='solaris.smf.value.routing' />
	</property_group>

	<template>
		<common_name>
			<loctext xml:lang='C'>
				in.ndpd IPv6 neighbour discovery daemon
			</loctext>
		</common_name>
		<documentation>
			<manpage title='in.ndpd' section='8'
				manpath='/usr/share/man' />
		</documentation>
	</template>
	</instance>
	<stability value='Unstable' />
</service>

</service_bundle>
