# usr/src/cmd/power/power.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
 Copyright 2005 Sun Microsystems, Inc.  All rights reserved.
 Use is subject to license terms.

 CDDL HEADER START

 The contents of this file are subject to the terms of the
 Common Development and Distribution License, Version 1.0 only
 (the "License").  You may not use this file except in compliance
 with the License.

 You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
 or http://www.opensolaris.org/os/licensing.
 See the License for the specific language governing permissions
 and limitations under the License.

 When distributing Covered Code, include this CDDL HEADER in each
 file and include the License file at usr/src/OPENSOLARIS.LICENSE.
 If applicable, add the following below this CDDL HEADER, with the
 fields enclosed by brackets "[]" replaced with your own identifying
 information: Portions Copyright [yyyy] [name of copyright owner]

 CDDL HEADER END

    NOTE:  This service manifest is not editable; its contents will
    be overwritten by package or patch operations, including
    operating system upgrade.  Make customizations in a different
    file.
-->

<service_bundle type='manifest' name='SUNWpmr:power'>

<service
	name='system/power'
	type='service'
	version='1'>

	<create_default_instance enabled='false' />

	<single_instance />

	<dependency
		name='usr'
		type='service'
		grouping='require_all'
		restart_on='none'>
		<service_fmri value='svc:/system/filesystem/minimal' />
	</dependency>

	<dependent
		name='power_multi-user'
		grouping='optional_all'
		restart_on='none'>
			<service_fmri value='svc:/milestone/multi-user' />
	</dependent>

	<exec_method
		type='method'
		name='start'
		exec='/lib/svc/method/svc-power %m'
		timeout_seconds='60'>
		<method_context working_directory='/'>
			<method_credential user='root' group='root'/>
		</method_context>
	</exec_method>

	<exec_method
		type='method'
		name='stop'
		exec='/lib/svc/method/svc-power %m'
		timeout_seconds='60'>
		<method_context working_directory='/'>
			<method_credential user='root' group='root'/>
		</method_context>
	</exec_method>

	<property_group name='startd' type='framework'>
		<propval name='duration' type='astring'
			value='transient' />
	</property_group>

	<property_group name='general' type='framework'>
		<!-- to start stop power daemon -->
		<propval name='action_authorization' type='astring'
			value='solaris.smf.manage.power' />
	</property_group>


	<stability value='Unstable' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
			power management
			</loctext>
		</common_name>
		<documentation>
			<manpage title='powerd' section='8' manpath='/usr/share/man' />
			<manpage title='pmconfig' section='8' manpath='/usr/share/man' />
			<manpage title='dtpower' section='8' manpath='/usr/share/man' />
			<manpage title='power.conf' section='5' manpath='/usr/share/man' />
		</documentation>
	</template>
</service>

</service_bundle>
