# usr/src/cmd/rcap/rcapd/rcap.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
 Copyright 2007 Sun Microsystems, Inc.  All rights reserved.
 Use is subject to license terms.

 CDDL HEADER START

 The contents of this file are subject to the terms of the
 Common Development and Distribution License (the "License").
 You may not use this file except in compliance with the License.

 You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
 or http://www.opensolaris.org/os/licensing.
 See the License for the specific language governing permissions
 and limitations under the License.

 When distributing Covered Code, include this CDDL HEADER in each
 file and include the License file at usr/src/OPENSOLARIS.LICENSE.
 If applicable, add the following below this CDDL HEADER, with the
 fields enclosed by brackets "[]" replaced with your own identifying
 information: Portions Copyright [yyyy] [name of copyright owner]

 CDDL HEADER END

    NOTE:  This service manifest is not editable; its contents will
    be overwritten by package or patch operations, including
    operating system upgrade.  Make customizations in a different
    file.
-->

<service_bundle type='manifest' name='SUNWrcapr:rcap'>

<service
	name='system/rcap'
	type='service'
	version='1'>

	<single_instance />

	<dependency
		name='usr'
		type='service'
		grouping='require_all'
		restart_on='none'>
		<service_fmri value='svc:/system/filesystem/minimal' />
	</dependency>

	<dependency
		name='manifest'
		type='service'
		grouping='optional_all'
		restart_on='none'>
			<service_fmri value='svc:/system/manifest-import' />
	</dependency>

	<dependent
		name='rcap_multi-user'
		grouping='optional_all'
		restart_on='none'>
			<service_fmri value='svc:/milestone/multi-user' />
	</dependent>

	<instance name='default' enabled='false'>

	<exec_method
		type='method'
		name='start'
		exec='/usr/lib/rcap/rcapd'
		timeout_seconds='60'>
		<method_context working_directory='/'>
			<method_credential user='root'/>
		</method_context>
	</exec_method>

	<exec_method
		type='method'
		name='refresh'
		exec=':kill -HUP'
		timeout_seconds='60' />

	<exec_method
		type='method'
		name='stop'
		exec=':kill'
		timeout_seconds='60' />

	<property_group name='config' type='application'>
		<propval name='pressure' type='count' value='0' />
		<propval name='reconfig_interval' type='count' value='60' />
		<propval name='walk_interval' type='count' value='15' />
		<propval name='report_interval' type='count' value='5' />
		<propval name='rss_sample_interval' type='count' value='5' />
        </property_group>

	</instance>

	<stability value='Unstable' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
			resource capping daemon
			</loctext>
		</common_name>
		<documentation>
			<manpage title='rcapd' section='8'
			     manpath='/usr/share/man' />
			<manpage title='rcapstat' section='8'
			     manpath='/usr/share/man' />
			<manpage title='rcapadm' section='8'
			     manpath='/usr/share/man' />
		</documentation>
	</template>
</service>

</service_bundle>
