# usr/src/cmd/cmd-inet/usr.lib/in.discardd/discard.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>

<!--
 Copyright 2005 Sun Microsystems, Inc.  All rights reserved.
 Use is subject to license terms.

 CDDL HEADER START

 The contents of this file are subject to the terms of the
 Common Development and Distribution License, Version 1.0 only
 (the "License").  You may not use this file except in compliance
 with the License.

 You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
 or http://www.opensolaris.org/os/licensing.
 See the License for the specific language governing permissions
 and limitations under the License.

 When distributing Covered Code, include this CDDL HEADER in each
 file and include the License file at usr/src/OPENSOLARIS.LICENSE.
 If applicable, add the following below this CDDL HEADER, with the
 fields enclosed by brackets "[]" replaced with your own identifying
 information: Portions Copyright [yyyy] [name of copyright owner]

 CDDL HEADER END

	NOTE:  This service manifest is not editable; its contents will
	be overwritten by package or patch operations, including
	operating system upgrade.  Make customizations in a different
	file.

	Service manifests for discard
-->

<service_bundle type='manifest' name='SUNWcnsr:discard'>

<service name='network/discard' type='service' version='1'>

	<restarter>
	    <service_fmri value='svc:/network/inetd:default' />
	</restarter>

	<property_group name='inetd' type='framework'>
	    <stability value='Evolving' />
	    <propval name='name' type='astring' value='discard' />
	    <propval name='isrpc' type='boolean' value='false' />
	</property_group>

	<instance name='dgram' enabled='false' >
	    <exec_method
		type='method'
		name='inetd_start'
		exec='/usr/lib/inet/in.discardd -d'
		timeout_seconds='0'>
		    <method_context working_directory='/'>
			<method_credential user='root' group='root' />
		    </method_context>
	    </exec_method>

	    <exec_method
		type='method'
		name='inetd_offline'
	        exec=':kill_process'
		timeout_seconds='0'>
	    </exec_method>

	    <exec_method
		type='method'
		name='inetd_disable'
		exec=':kill'
		timeout_seconds='0'>
	    </exec_method>

		<property_group name='inetd' type='framework'>
		    <propval name='endpoint_type' type='astring'
			    value='dgram' />
		    <propval name='proto' type='astring' value='udp6' />
		    <propval name='wait' type='boolean' value='true' />
		</property_group>
	</instance>

	<instance name='stream' enabled='false' >
	    <exec_method
		type='method'
		name='inetd_start'
		exec='/usr/lib/inet/in.discardd -s'
		timeout_seconds='0'>
		    <method_context working_directory='/'>
			<method_credential user='root' group='root' />
		    </method_context>
	    </exec_method>

	    <exec_method
		type='method'
		name='inetd_disable'
		exec=':kill'
	        timeout_seconds='0'>
	    </exec_method>

	    <property_group name='inetd' type='framework'>
		<propval name='endpoint_type' type='astring'
			    value='stream' />
		<propval name='proto' type='astring' value='tcp6' />
		<propval name='wait' type='boolean' value='false' />
	    </property_group>
	</instance>

	<stability value='Evolving' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
			discard
			</loctext>
		</common_name>
		<documentation>
		</documentation>
	</template>
</service>

</service_bundle>
