# usr/src/cmd/print/bsd-sysv-commands/rfc1179.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>

<!--
 CDDL HEADER START

 The contents of this file are subject to the terms of the
 Common Development and Distribution License (the "License").
 You may not use this file except in compliance with the License.

 You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
 or http://www.opensolaris.org/os/licensing.
 See the License for the specific language governing permissions
 and limitations under the License.

 When distributing Covered Code, include this CDDL HEADER in each
 file and include the License file at usr/src/OPENSOLARIS.LICENSE.
 If applicable, add the following below this CDDL HEADER, with the
 fields enclosed by brackets "[]" replaced with your own identifying
 information: Portions Copyright [yyyy] [name of copyright owner]

 CDDL HEADER END

 Copyright 2009 Sun Microsystems, Inc.  All rights reserved.
 Use is subject to license terms.

	NOTE:  This service manifest is not editable; its contents will
	be overwritten by package or patch operations, including
	operating system upgrade.  Make customizations in a different
	file.
-->

<service_bundle type='manifest' name='SUNWpsr:rfc1179'>

<service
    name='application/print/rfc1179'
    type='service'
    version='1'>

	<create_default_instance enabled='false' />

	<restarter>
		<service_fmri value='svc:/network/inetd:default' />
	</restarter>

	<dependency
		name='lpsched'
		grouping='require_all'
		restart_on='refresh'
		type='service'>
		<service_fmri value='svc:/application/print/server' />
	</dependency>

	<exec_method
		type='method'
		name='inetd_start'
		exec='/usr/lib/print/in.lpd -u lp -d /var/run/in.lpd'
		timeout_seconds='0'>
		<method_context working_directory='/'>
			<method_credential user='root' group='lp' />
		</method_context>
	</exec_method>


	<exec_method
		type='method'
		name='inetd_disable'
		exec=':kill'
		timeout_seconds='0'>
	</exec_method>

	<property_group name='inetd' type='framework'>
		<stability value='Evolving' />
		<propval name='endpoint_type' type='astring' value='stream' />
		<propval name='name' type='astring' value='printer' />
		<propval name='wait' type='boolean' value='false' />
		<propval name='isrpc' type='boolean' value='false' />
		<propval name='proto' type='astring' value='tcp6' />
	</property_group>

	<property_group name='general' type='framework'>
		<!-- to start/stop rfc1179 listening -->
		<propval name='action_authorization' type='astring'
			value='solaris.print.admin' />
		<propval name='value_authorization' type='astring'
			value='solaris.print.admin' />
	</property_group>

	<property_group name='firewall_context' type='com.sun,fw_definition'>
		<propval name='ipf_method' type='astring'
			value='/lib/svc/method/print-svc ipfilter svc:/application/print/server:default' />
	</property_group>

	<stability value='Unstable' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
			BSD print protocol adapter
			</loctext>
		</common_name>
		<documentation>
			<manpage title='in.lpd' section='8'
			    manpath='/usr/share/man' />
		</documentation>
	</template>

</service>

</service_bundle>
