# usr/src/cmd/svc/profile/inetd_generic.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!--
    CDDL HEADER START
   
    The contents of this file are subject to the terms of the
    Common Development and Distribution License (the "License").
    You may not use this file except in compliance with the License.
   
    You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
    or http://www.opensolaris.org/os/licensing.
    See the License for the specific language governing permissions
    and limitations under the License.
   
    When distributing Covered Code, include this CDDL HEADER in each
    file and include the License file at usr/src/OPENSOLARIS.LICENSE.
    If applicable, add the following below this CDDL HEADER, with the
    fields enclosed by brackets "[]" replaced with your own identifying
    information: Portions Copyright [yyyy] [name of copyright owner]
   
    CDDL HEADER END
   
    Copyright 2010 Sun Microsystems, Inc.  All rights reserved.
    Use is subject to license terms.

    Service profile to enable default inetd services.

    NOTE:  Service profiles delivered by this package are not editable,
    and their contents will be overwritten by package or patch
    operations, including operating system upgrade.  Make customizations
    in a distinct file.
-->
<service_bundle type='profile' name='default'>
	<!--
		inetd(1M) services
	-->
	<service name='network/finger' version='1' type='service'>
		<instance name='default' enabled='true'/>
	</service>
	<service name='network/login' version='1' type='service'>
		<instance name='rlogin' enabled='true'/>
	</service>
	<service name='network/shell' version='1' type='service'>
		<instance name='default' enabled='true'/>
	</service>
	<service name='network/telnet' version='1' type='service'>
		<instance name='default' enabled='true'/>
	</service>

	<!--
		inetd(1M) RPC services
	-->
	<service name='network/nfs/rquota' version='1' type='service'>
		<instance name='default' enabled='true'/>
	</service>
	<service name='network/rpc/gss' version='1' type='service'>
		<instance name='default' enabled='true'/>
	</service>
	<service name='network/rpc/mdcomm' version='1' type='service'>
		<instance name='default' enabled='true'/>
	</service>
	<service name='network/rpc/meta' version='1' type='service'>
		<instance name='default' enabled='true'/>
	</service>
	<service name='network/rpc/metamed' version='1' type='service'>
		<instance name='default' enabled='true'/>
	</service>
	<service name='network/rpc/metamh' version='1' type='service'>
		<instance name='default' enabled='true'/>
	</service>
	<service name='network/rpc/rstat' version='1' type='service'>
		<instance name='default' enabled='true'/>
	</service>
	<service name='network/rpc/rusers' version='1' type='service'>
		<instance name='default' enabled='true'/>
	</service>
	<service name='network/rpc/smserver' version='1' type='service'>
		<instance name='default' enabled='true'/>
	</service>
	<service name='network/security/ktkt_warn' version='1' type='service'>
		<instance name='default' enabled='true'/>
	</service>

	<!--
		Enable CDE services.
	-->
	<service name='network/rpc/cde-ttdbserver' version='1' type='service'>
		<instance name='tcp' enabled='true' />
	</service>
	<service name='network/rpc/cde-calendar-manager' version='1'
		type='service'>
		<instance name='default' enabled='true' />
	</service>

</service_bundle>
