# usr/src/cmd/svc/profile/ns_none.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!--
 Copyright 2009 Sun Microsystems, Inc.  All rights reserved.
 Use is subject to license terms.

 CDDL HEADER START

 The contents of this file are subject to the terms of the
 Common Development and Distribution License (the "License").
 You may not use this file except in compliance with the License.

 You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
 or http://www.opensolaris.org/os/licensing.
 See the License for the specific language governing permissions
 and limitations under the License.

 When distributing Covered Code, include this CDDL HEADER in each
 file and include the License file at usr/src/OPENSOLARIS.LICENSE.
 If applicable, add the following below this CDDL HEADER, with the
 fields enclosed by brackets "[]" replaced with your own identifying
 information: Portions Copyright [yyyy] [name of copyright owner]

 CDDL HEADER END

    NOTE:  This service profile is not editable; its contents will be
    overwritten by package or patch operations, including operating
    system upgrade.

    Service profile to deactivate all network-dependent name services.
-->
<service_bundle type='profile' name='default'>
	<!-- DNS -->
	<service name='network/dns/client' version='1' type='service'>
		<instance name='default' enabled='false'/>
	</service>
	<!-- LDAP -->
	<service name='network/ldap/client' version='1' type='service'>
		<instance name='default' enabled='false'/>
	</service>
	<!-- NIS, client and server -->
	<service name='network/nis/client' version='1' type='service'>
		<instance name='default' enabled='false'/>
	</service>
	<service name='network/nis/server' version='1' type='service'>
		<instance name='default' enabled='false'/>
	</service>
	<service name='network/nis/passwd' version='1' type='service'>
		<instance name='default' enabled='false'/>
	</service>
	<service name='network/nis/update' version='1' type='service'>
		<instance name='default' enabled='false'/>
	</service>
	<service name='network/nis/xfr' version='1' type='service'>
		<instance name='default' enabled='false'/>
	</service>
	<!-- supporting services for NIS -->
	<service name='network/rpc/keyserv' version='1' type='service'>
		<instance name='default' enabled='false'/>
	</service>
</service_bundle>
