# usr/src/cmd/cmd-inet/usr.sbin/ipsecutils/tcpkey.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
 Copyright 2007 Sun Microsystems, Inc.  All rights reserved.
 Use is subject to license terms.

 CDDL HEADER START

 The contents of this file are subject to the terms of the
 Common Development and Distribution License (the "License").
 You may not use this file except in compliance with the License.

 You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
 or http://www.opensolaris.org/os/licensing.
 See the License for the specific language governing permissions
 and limitations under the License.

 When distributing Covered Code, include this CDDL HEADER in each
 file and include the License file at usr/src/OPENSOLARIS.LICENSE.
 If applicable, add the following below this CDDL HEADER, with the
 fields enclosed by brackets "[]" replaced with your own identifying
 information: Portions Copyright [yyyy] [name of copyright owner]

 Copyright 2024 Oxide Computer Company

 CDDL HEADER END

	NOTE:  This service manifest is not editable; its contents will
	be overwritten by package or patch operations, including
	operating system upgrade.  Make customizations in a different
	file.
-->
<service_bundle type='manifest' name='SUNWcsr:tcpkey'>

<service
	name='network/tcpkey'
	type='service'
	version='1'>

	<!-- The 'tcpkey' service is delivered disabled
	because there is not a default configuration file.
	See note below on changing the default configuration file. -->

	<create_default_instance enabled='false' />

	<single_instance />

	<!-- Read/Write access to /var/run required for lock files -->
	<dependency
		name='filesystem'
		grouping='require_all'
		restart_on='none'
		type='service'>
		<service_fmri
			value='svc:/system/filesystem/minimal'
		/>
	</dependency>

	<!-- If we are enabled, we should be running fairly early -->
	<dependent
		name='tcpkey-network'
		grouping='optional_all'
		restart_on='none'>
		<service_fmri
			value='svc:/milestone/network'
		/>
	</dependent>

	<method_context working_directory='/'>
		<method_credential user="root"
			group="root"
			privileges="basic_test,file_read,proc_fork,proc_exec,sys_ip_config"
			/>
	</method_context>

	<exec_method
		type='method'
		name='start'
		exec='/lib/svc/method/net-tcpkey %{config/config_file}'
		timeout_seconds='60'
	/>

	<!-- To prevent tcpkey generating warnings about duplicate
	SAs when the service is refreshed, tcpkey will flush the
	existing SAs when it's called from smf(7). -->

	<exec_method
		type='method'
		name='refresh'
		exec='/lib/svc/method/net-tcpkey %{config/config_file}'
		timeout_seconds='60'
	/>

	<exec_method
		type='method'
		name='stop'
		exec='/usr/sbin/tcpkey flush'
		timeout_seconds='60'
	/>

	<property_group name='general' type='framework'>
		<!-- A user with this authorization can:

			svcadm restart tcpkey
			svcadm refresh tcpkey
			svcadm mark <state> tcpkey
			svcadm clear tcpkey

		see auths(1) and user_attr(5)-->

		<propval
			name='action_authorization'
			type='astring'
			value='solaris.smf.manage.tcpkey'
		/>
		<!-- A user with this authorization can:

			svcadm disable tcpkey
			svcadm enable tcpkey

		see auths(1) and user_attr(5)-->

		<propval
			name='value_authorization'
			type='astring'
			value='solaris.smf.manage.tcpkey'
		/>
	</property_group>

	<!-- The properties defined below can be changed by a user
	with 'solaris.smf.value.tcpkey' authorization using the
	svccfg(8) command.

	E.g.:

	svccfg -s tcpkey setprop config/config_file = /new/config_file

	The new configurations will be read on service refresh:

	svcadm refresh tcpkey

	Note: svcadm disable/enable does not use the new property
	until after the service has been refreshed.

	***Do not edit this manifest to change these properties! -->

	<property_group name='config' type='application'>
		<propval
			name='config_file'
			type='astring'
			value='/etc/inet/secret/tcpkeys'
		/>
		<propval
			name='value_authorization'
			type='astring'
			value='solaris.smf.value.tcpkey'
		/>
	</property_group>

	<property_group name='startd' type='framework'>
		<propval
			name='duration'
			type='astring'
			value='transient'
		/>
	</property_group>

	<stability value='Unstable' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
				TCP SADB startup
			</loctext>
		</common_name>
		<description>
			<loctext xml:lang='C'>
				Loads static security associations
			</loctext>
		</description>
		<documentation>
			<manpage title='tcpkey' section='8'
				manpath='/usr/share/man' />
		</documentation>
	</template>
</service>
</service_bundle>

