# usr/src/lib/libresolv2/install.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
 Copyright (c) 2010, Oracle and/or its affiliates. All rights reserved.

 CDDL HEADER START

 The contents of this file are subject to the terms of the
 Common Development and Distribution License (the "License").
 You may not use this file except in compliance with the License.

 You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
 or http://www.opensolaris.org/os/licensing.
 See the License for the specific language governing permissions
 and limitations under the License.

 When distributing Covered Code, include this CDDL HEADER in each
 file and include the License file at usr/src/OPENSOLARIS.LICENSE.
 If applicable, add the following below this CDDL HEADER, with the
 fields enclosed by brackets "[]" replaced with your own identifying
 information: Portions Copyright [yyyy] [name of copyright owner]

 CDDL HEADER END

	NOTE:  This service manifest is not editable; its contents will
	be overwritten by package or patch operations, including
	operating system upgrade.  Make customizations in a different
	file.
-->

<service_bundle type='manifest' name='SUNWcsr:dns-install'>

<service
	name='network/dns/install'
	type='service'
	version='1'>

	<create_default_instance enabled='false' />

	<single_instance />

	<dependency
		name='filesystem'
		grouping='require_all'
		restart_on='none'
		type='service'>
		<service_fmri value='svc:/system/filesystem/root' />
		<service_fmri value='svc:/system/filesystem/usr' />
		<service_fmri value='svc:/system/filesystem/minimal' />
	</dependency>

	<dependent name='network-service'
		grouping='optional_all'
		restart_on='none'>
		<service_fmri value='svc:/network/service' />
	</dependent>

	<dependent name='dns-client'
		grouping='optional_all'
		restart_on='none'>
		<service_fmri value='svc:/network/dns/client' />
	</dependent>

	<exec_method
		type='method'
		name='start'
		exec='/lib/svc/method/dns-install'
		timeout_seconds='30' />

	<exec_method
		type='method'
		name='stop'
		exec=':true'
		timeout_seconds='0' />

	<property_group name='startd' type='framework'>
		<propval name='duration' type='astring' value='transient' />
	</property_group>

	<property_group name='install_props' type='application'>
		<property name='nameserver' type='net_address'>
			<net_address_list>
				<value_node value='0.0.0.0' />
			</net_address_list>
		</property>
		<propval name='domain' type='astring' value='' />
		<property name='search' type='astring'>
			<astring_list>
				<value_node value='' />
			</astring_list>
		</property>
	</property_group>

	<stability value='Unstable' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
				DNS resolver install
			</loctext>
		</common_name>
		<documentation>
			<manpage title='resolver' section='3RESOLV'
				manpath='/usr/share/man' />
		</documentation>
		<pg_pattern name='install_props' type='application'
		    target='this' required='false'>
			<description> <loctext xml:lang='C'>
			    Install derived configuration data used to configure an initial DNS client.
			 </loctext> </description>
			<prop_pattern name='nameserver' type='net_address'
			    required='true'>
				<description> <loctext xml:lang='C'>
				    The value used to construct the "nameserver" directive in resolv.conf(5).
				</loctext> </description>
				<cardinality min='1' max='3'/>
			</prop_pattern>
			<prop_pattern name='domain' type='astring'
			    required='false'>
				<description> <loctext xml:lang='C'>
				    The value used to construct the "domain" directive in resolv.conf(5).
				</loctext> </description>
				<cardinality min='1' max='1'/>
			</prop_pattern>
			<prop_pattern name='search' type='astring'
			    required='false'>
				<description> <loctext xml:lang='C'>
				    The value used to construct the "search" directive in resolv.conf(5).
				</loctext> </description>
				<cardinality min='1' max='6'/>
			</prop_pattern>
		</pg_pattern>
	</template>

</service>

</service_bundle>
